# root/etc/services.d/sonos-web/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc cd /root/.sonos-web/ node src/server.js
